# lib/cron.d/m2backup
# cron-jobs for mark2-backup
# https://github.com/yeahwhat-mc/mark2-backup

0 3 * * *        minecraft       /usr/bin/m2backup backup main
0 3 * * *        minecraft       /usr/bin/m2backup clear main
